# crontab
* * * * * /opt/adjust-fans.sh >> /var/log/temp.log
